# temper.service
[Unit]
Description=Read and publish temperature

[Service]
Type=simple
ExecStart=/opt/temperature_exporter/temp.sh
